# net-im/palaver/files/palaver-initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use jabber-server
}

start() {
	ebegin "Starting Palaver MUC service"
	start-stop-daemon --start --pidfile /var/run/jabber/palaver.pid -u jabber -g jabber \
		--exec /usr/bin/twistd -- -l /var/log/jabber/palaver.log --pidfile /var/run/jabber/palaver.pid \
		palaver -c ${CONFIG}
	eend $?
}

stop() {
	ebegin "Stopping Palaver MUC service"
	start-stop-daemon --stop --quiet --pidfile /var/run/jabber/palaver.pid
	eend $?
}
