# net-im/punjab/files/punjab-initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use jabber-server
}

start() {
	ebegin "Starting Punjab BOSH connection manager"
	start-stop-daemon --start --pidfile /var/run/jabber/punjab.pid -u jabber -g jabber \
		--exec /usr/bin/twistd -- -l /var/log/jabber/punjab.log --pidfile /var/run/jabber/punjab.pid \
		-y /etc/jabber/punjab.tac
# -c ${CONFIG}
	eend $?
}

stop() {
	ebegin "Stopping Punjab BOSH connection manager"
	start-stop-daemon --stop --quiet --pidfile /var/run/jabber/punjab.pid
	eend $?
}
