# mail-filter/policyd/files/policyd.service
[Unit]
Description=Policyd

[Service]
ExecStart=/usr/sbin/policyd -c /etc/policyd.conf

[Install]
WantedBy=multi-user.target
